# shutdown-if-inactive.crontab
*/15 * * * * /home/ubuntu/pcw/shutdown-if-inactive.sh >> /var/log/shutdown-if-inactive.log 2>&1